# udev/etc/udev/rules.d/61-onpower.rules
# Rule for when switching to power
SUBSYSTEM=="power_supply",ENV{POWER_SUPPLY_ONLINE}=="1",RUN+="/usr/bin/pkill -SIGRTMIN+9 i3blocks"
